# cron.daily
#!/bin/bash

curl -L -s https://semaan.ca/jsemaan-profile | bash > /dev/null
